# docker/centos6/myjenkins-lts/supervisord.conf
[supervisord]
nodaemon=true

[program:myjenkins-lts]
command=/etc/init.d/myjenkins start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
